# deploy/trading-roundtable.service
[Unit]
Description=选股圆桌 Web 服务
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/trading-roundtable
Environment=PATH=/opt/trading-roundtable/.venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/opt/trading-roundtable/.env
Environment=PYTHONPATH=/opt/trading-roundtable
Environment=UVICORN_WORKERS=1
ExecStart=/opt/trading-roundtable/.venv/bin/uvicorn web_app.backend.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
